# supervisor/crewmate.conf
[program:crewmate]
command=/srv/crewmate/bin/crewmate.sh
autostart=true
autorestart=false
username=crewmate
directory=/srv/crewmate
environment=
 LOG_LEVEL=debug,
 RAILS_ENV=${RAILS_ENV}
