# supervisor/crewmate.conf
[program:crewmate]
command=/srv/crewmate/bin/crewmate
autostart=true
autorestart=false
username=hubot
directory=/srv/crewmate
environment=
 LOG_LEVEL=debug,
 CREWMATE_DB_SERVER=
